# systemd/opcuator-github-sync.timer
[Unit]
Description=Run OPCUAtor GitHub sync every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
