# systemd/opcuator-github-sync.service
[Unit]
Description=Sync OPCUAtor repository to GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/opcuator
EnvironmentFile=-/etc/opcuator/github-sync.env
ExecStart=/opt/opcuator/scripts/opcuator-github-sync.sh
